# deploy/oci/zipmould-viz.service
[Unit]
Description=ZipMould FastAPI backend
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=zipmould
Group=zipmould
WorkingDirectory=/opt/zipmould/slime-mould
EnvironmentFile=/etc/zipmould-viz.env
ExecStart=/home/zipmould/.local/bin/uv run zipmould viz serve --host 127.0.0.1 --port 8000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
